# agent/deploy/nexus-agent.service
[Unit]
Description=Nexus Agent - Infrastructure Management
Documentation=https://github.com/nexus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nexus-agent
Group=nexus-agent
# systemd cree /var/lib/nexus-agent automatiquement (owner=nexus-agent, mode=0700)
StateDirectory=nexus-agent
StateDirectoryMode=0700
ExecStart=/usr/local/bin/nexus-agent
EnvironmentFile=/etc/nexus/agent.env
Restart=always
RestartSec=10
WatchdogSec=120

# ===================== Sécurité =====================

# L'agent tourne sous nexus-agent (non-root)
# Les commandes privilégiées passent par sudo (voir /etc/sudoers.d/nexus-agent)
ProtectHome=true
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
# RestrictSUIDSGID retire : sudo est un binaire SUID necessaire pour les actions privilegiees
# La protection reste via le sudoers ciblé (whitelist + NOEXEC)
RestrictRealtime=true
LockPersonality=true

# Capabilities Linux minimales pour le monitoring
# Ambient: promues au non-root nexus-agent
# Pas de CapabilityBoundingSet : apt/sudo ont besoin du full set (chown, fowner, etc.)
AmbientCapabilities=CAP_NET_RAW CAP_SYS_PTRACE CAP_DAC_READ_SEARCH

# SystemCallFilter retire : interfere avec sudo (audit syscalls) et apt-get
# La protection reste via les autres directives (Protect*, capabilities, sudoers)
SystemCallArchitectures=native

# Répertoires accessibles
ReadOnlyPaths=/proc /sys /etc/os-release
ReadWritePaths=/var/lib/nexus/keys /var/log/nexus

# Limites de ressources
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
