# agent/deploy/nexus-agent.service
[Unit]
Description=Nexus Agent - Infrastructure Management
Documentation=https://github.com/nexus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nexus-agent
Group=nexus-agent
# Creer le repertoire avant ReadWritePaths (evite l'erreur namespace si absent)
ExecStartPre=+/bin/mkdir -p /var/tmp/nexus-agent
ExecStartPre=+/bin/chown nexus-agent:nexus-agent /var/tmp/nexus-agent
ExecStartPre=+/bin/chmod 0700 /var/tmp/nexus-agent
ExecStart=/usr/local/bin/nexus-agent
EnvironmentFile=/etc/nexus/agent.env
Restart=always
RestartSec=10
WatchdogSec=120

# ===================== Sécurité =====================

# L'agent tourne sous nexus-agent (non-root)
# Les commandes privilégiées passent par sudo (voir /etc/sudoers.d/nexus-agent)
ProtectHome=true
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true

# Capabilities Linux minimales pour le monitoring
AmbientCapabilities=CAP_NET_RAW CAP_SYS_PTRACE CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_NET_RAW CAP_SYS_PTRACE CAP_DAC_READ_SEARCH

# Limiter les syscalls aux essentiels
SystemCallFilter=@system-service @file-system @network-io @process
SystemCallArchitectures=native

# Répertoires accessibles
ReadOnlyPaths=/proc /sys /etc/os-release
ReadWritePaths=/var/lib/nexus/keys /var/log/nexus /var/tmp/nexus-agent

# Limites de ressources
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
